# systemd/block-fw-update.service
[Unit]
Description=block-fw ipset update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/block-fw update
